# cron/crontab
#crontab
@reboot root *default/initialize
